# assets/etc/supervisor/conf.d/nginx.ini
[program:nginx]
username=www-data

command=/usr/sbin/nginx -g "daemon off;"

autostart=true
autorestart=true
stopasgroup=true
