# assets/etc/supervisor/conf.d/postal.ini
[program:postal]
username=postal

command=sudo -i -u postal /usr/bin/postal run

autostart=true
autorestart=true
stopasgroup=true
